# debian/etc/cron.d/pbs_licenses
*/1 * * * * root /var/spool/pbs/clean_licenses.py &>/dev/null
*/1 * * * * root /var/spool/pbs/licenses/update_cache_all.sh
